# srcs/requirements/nginx/tools/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=rc-service nginx start
stdout_logfile=/var/log/supervisord/nginx.log
startsecs=0
