# srcs/requirements/mariadb/srcs/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord/supervisord.log
pidfile=/var/run/supervisord.pid

[program:mysql]
command=rc-service mariadb start
stdout_logfile=/var/log/supervisord/mysql.log
